# supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/run-worker.sh
