# supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/run-worker.sh

[program:ssh]
command=/usr/sbin/sshd -D